# packer/userWebApp.service
[Unit]
Description=webapp
After=syslog.target

[Service]
User=ec2-user
EnvironmentFile=/etc/systemd/system/service.env
SuccessExitStatus=143
Restart=always
RestartSec=10
SyslogIdentifier=Inventory-webapp
StandardOutput=syslog
StandardError=syslog

ExecStart=/usr/bin/java -jar /usr/bin/UserWebApp-0.0.1-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target



